# var/lib/opkg/info/gadget-init.postrm
#!/bin/sh
update-rc.d $D storage-gadget-init remove
