# var/lib/opkg/info/libgles-omap3.list
/usr/lib/libEGL.so
/usr/lib/libsrv_um.so.1.9.2188537
/usr/lib/libIMGegl.so.1.9
/usr/lib/libglslcompiler_eglimage.so.1.9.2188537
/usr/lib/libEGL_eglimage.so
/usr/lib/libGLESv2_eglimage.so.1.9.2188537
/usr/lib/libPVRScopeServices.so
/usr/lib/libsrv_um_dri.so.1
/usr/lib/libpvrEWS_WSEGL.so.1
/usr/lib/libsrv_init.so
/usr/lib/libsrv_init.so.1.9
/usr/lib/libglslcompiler_eglimage.so.1
/usr/lib/libsrv_init.so.1
/usr/lib/libGLESv2_eglimage.so.1.9
/usr/lib/libIMGegl.so
/usr/lib/libpvr2d.so
/usr/lib/libpvr2d.so.1.9
/usr/bin/SGX/trainingcourses/Raw/
/usr/lib/libpvr2d.so.1
/usr/lib/libGLES_CM_eglimage.so.1.9
/usr/lib/libGLESv2_eglimage.so.1
/usr/lib/libpvrPVR2D_BLITWSEGL.so
/usr/lib/libsrv_init.so.1.9.2188537
/usr/lib/libGLESv2.so.1.9
/usr/lib/libEGL.so.1.9
/etc/init.d/pvr-init
/usr/lib/libusc.so
/usr/lib/libGLESv2_eglimage.so
/usr/bin/gles1_texture_stream
/usr/lib/libpvrPVR2D_DRIWSEGL.so.1
/usr/lib/libpvrEWS_WSEGL.so.1.9.2188537
/usr/lib/libEGL_eglimage.so.1
/usr/lib/libglslcompiler.so.1.9.2188537
/usr/lib/libpvr2d.so.1.9.2188537
/etc/powervr.ini
/usr/lib/libGLES_CM.so.1.9
/usr/lib/libsrv_um.so.1.9
/usr/lib/libIMGegl_eglimage.so
/usr/lib/libpvrPVR2D_FRONTWSEGL.so
/usr/lib/libpvrEWS_WSEGL.so
/usr/lib/libglslcompiler.so.1.9
/usr/lib/libIMGegl_eglimage.so.1
/usr/lib/libGLES_CM_eglimage.so
/usr/lib/libpvrPVR2D_LINUXFBWSEGL.so
/usr/lib/libIMGegl_eglimage.so.1.9
/usr/lib/libsrv_um.so.1
/usr/lib/libGLES_CM.so.1.9.2188537
/usr/lib/libglslcompiler_eglimage.so
/usr/lib/libglslcompiler.so.1
/usr/lib/libglslcompiler.so
/usr/lib/libEGL_eglimage.so.1.9
/usr/lib/libsrv_um.so
/usr/lib/libpvrPVR2D_FLIPWSEGL.so
/usr/lib/libGLESv2.so
/usr/lib/libIMGegl.so.1
/usr/lib/libGLES_CM.so.1
/usr/lib/libglslcompiler_eglimage.so.1.9
/usr/lib/libEGL.so.1
/usr/lib/libGLESv2.so.1.9.2188537
/usr/lib/libpvrEWS_WSEGL.so.1.9
/usr/bin/cputype
/usr/bin/ES2.0/
/usr/lib/libEGL_eglimage.so.1.9.2188537
/usr/bin/SGX/trainingcourses/
/usr/lib/libews.so
/usr/bin/gles2_texture_stream
/usr/lib/libsrv_um_dri.so.1.9
/usr/lib/libpvrPVR2D_DRIWSEGL.so
/usr/lib/libpvrPVR2D_DRIWSEGL.so.1.9
/etc/udev/rules.d/99-bufferclass.rules
/usr/lib/libGLES_CM.so
/usr/lib/libGLES_CM_eglimage.so.1
/usr/lib/libGLES_CM_eglimage.so.1.9.2188537
/usr/lib/pvr_drv.so
/usr/bin/pvrsrvctl
/usr/lib/libIMGegl_eglimage.so.1.9.2188537
/usr/lib/libEGL.so.1.9.2188537
/usr/lib/libIMGegl.so.1.9.2188537
/usr/lib/libsrv_um_dri.so.1.9.2188537
/usr/lib/libpvrPVR2D_DRIWSEGL.so.1.9.2188537
/usr/lib/libGLESv2.so.1
