# var/lib/opkg/info/matrix-gui.postrm
#!/bin/sh
update-rc.d $D matrix-gui-2.0 remove
