# var/lib/opkg/info/mount-sdcard.postrm
#!/bin/sh
update-rc.d $D mount-sdcard remove
