# var/lib/opkg/info/parse-ip.postrm
#!/bin/sh
update-rc.d $D parse-ip remove
