# var/lib/opkg/info/wilink-calibrate.postrm
#!/bin/sh
update-rc.d $D calibrate.sh remove
